# lib/udev/rules.d/89-microcode.rules
KERNEL=="cpu0", ACTION=="add", RUN+="/sbin/modprobe microcode"
KERNEL=="microcode", ACTION=="add", RUN+="/sbin/microcode_ctl -Qu"
